# customisations/scripts/keepalive/keepalive.timer
[Unit]
Description=Run my keepalive script every 10 minutes

[Timer]
OnCalendar=*:0/4
Persistent=true

[Install]
WantedBy=timers.target
